# config/supervisor/conf.d/berrynet-darknet.conf
[program:darknet-service]
command=python3 /usr/lib/python3/dist-packages/berrynet/service/darknet_service.py --model_package tinyyolovoc --service_name detector --draw --debug
stdout_logfile=/var/log/berrynet/darknet-service-stdout.log
stdout_logfile_maxbytes=1048576
stderr_logfile=/var/log/berrynet/darknet-service-stderr.log
stderr_logfile_maxbytes=1048576
priority=10

[program:camera]
command=python3 /usr/lib/python3/dist-packages/berrynet/client/camera.py --fps 3
stdout_logfile=/var/log/berrynet/camera-stdout.log
stdout_logfile_maxbytes=1048576
stderr_logfile=/var/log/berrynet/camera-stderr.log
stderr_logfile_maxbytes=1048576
priority=30

[program:freeboard]
; Setup Freeboard for darknet-service
;   Step 1: Enter http://127.0.0.1:8080 in browser's URL bar, and press enter
;   Step 2: Click "LOAD FREEBOARD", and select BerryNet/config/dashboard-darknet.json
;   Step 3: Wait for seconds, you should see the inference result image and text on Freeboard
command=bash -c "cd /usr/lib/berrynet/dashboard/ && nodejs server.js"
stdout_logfile=/var/log/berrynet/freeboard-stdout.log
stdout_logfile_maxbytes=1048576
stderr_logfile=/var/log/berrynet/freeboard-stderr.log
stderr_logfile_maxbytes=1048576
priority=40
