# systemd/cleaner.timer
[Unit]
Description=Timer to clean snapshots every 5 min.

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
